# archive/system2/services/anvil-spy.service
[Unit]
Description=AnvilOS Spy (Audit Daemon)
After=network.target

[Service]
ExecStart=/usr/bin/python3 /mnt/anvil_temp/github/anvilos_alpha/runtime/agents/spy.py
Restart=always
User=aimeat
Environment=PYTHONPATH=/mnt/anvil_temp/github/anvilos_alpha/runtime

[Install]
WantedBy=multi-user.target
