# archive/system2/services/mainframe.initd
#!/sbin/openrc-run

name="Anvil Mainframe"
description="Mainframe Processor Daemon for AnvilOS"
command="/usr/bin/python3"
command_args="/mnt/anvil_temp/github/anvilos_alpha/runtime/processor_daemon.py"
pidfile="/run/mainframe.pid"
command_background="yes"
output_log="/var/log/mainframe.log"
error_log="/var/log/mainframe.log"

export PYTHONPATH="/usr/lib/python3.13/site-packages:/mnt/anvil_temp/github/anvilos_alpha/runtime"

depend() {
    need localmount
}
