# archive/oss_sovereignty/zfs-2.2.2/etc/systemd/system/zfs-volume-wait.service
[Unit]
Description=Wait for ZFS Volume (zvol) links in /dev
DefaultDependencies=no
After=systemd-udev-settle.service
After=zfs-import.target
ConditionPathIsDirectory=/sys/module/zfs

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/zfs
ExecStart=/usr/bin/zvol_wait

[Install]
WantedBy=zfs-volumes.target
