# oss_sovereignty/zfs-2.2.2/etc/init.d/zfs-import.in.mpy
{
  "module_name": "zfs-import.in",
  "hash_id": "90cfa06d6518025d6516f57bdcce4f078f9a60ad593262ce42b05bfe8803915c",
  "original_prompt": "Ingested from zfs-2.2.2/etc/init.d/zfs-import.in",
  "human_readable_source": "#!@DEFAULT_INIT_SHELL@\n# shellcheck disable=SC2154\n#\n# zfs-import    This script will import ZFS pools\n#\n# chkconfig:    2345 01 99\n# description:  This script will perform a verbatim import of ZFS pools\n#               during system boot.\n# probe: true\n#\n### BEGIN INIT INFO\n# Provides:          zfs-import\n# Required-Start:    mtab\n# Required-Stop:     $local_fs mtab\n# Default-Start:     S\n# Default-Stop:      0 1 6\n# X-Start-Before:    checkfs\n# X-Stop-After:      zfs-mount\n# Short-Description: Import ZFS pools\n# Description: Run the `zpool import` command.\n### END INIT INFO\n#\n# NOTE: Not having '$local_fs' on Required-Start but only on Required-Stop\n#       is on purpose. If we have '$local_fs' in both (and X-Start-Before=checkfs)\n#       we get conflicts - import needs to be started extremely early,\n#       but not stopped too late.\n#\n# Released under the 2-clause BSD license.\n#\n# This script is based on debian/zfsutils.zfs.init from the\n# Debian GNU/kFreeBSD zfsutils 8.1-3 package, written by Aurelien Jarno.\n\n# Source the common init script\n. @sysconfdir@/zfs/zfs-functions\n\n# ----------------------------------------------------\n\ndo_depend()\n{\n\tbefore swap\n\tafter sysfs udev\n\tkeyword -lxc -openvz -prefix -vserver\n}\n\n# Use the zpool cache file to import pools\ndo_verbatim_import()\n{\n\tif [ -f \"$ZPOOL_CACHE\" ]\n\tthen\n\t\tzfs_action \"Importing ZFS pool(s)\" \\\n\t\t\t\"$ZPOOL\" import -c \"$ZPOOL_CACHE\" -N -a\n\tfi\n}\n\n# Support function to get a list of all pools, separated with ';'\nfind_pools()\n{\n\tlocal pools\n\n\tpools=$(\"$@\" 2> /dev/null | \\\n\t\tsed -Ee '/pool:|^[a-zA-Z0-9]/!d' -e 's@.*: @@' | \\\n\t\tsort | \\\n\t\ttr '\\n' ';')\n\n\techo \"${pools%%;}\" # Return without the last ';'.\n}\n\n# Find and import all visible pools, even exported ones\ndo_import_all_visible()\n{\n\tlocal already_imported available_pools pool npools\n\tlocal exception dir ZPOOL_IMPORT_PATH RET=0 r=1\n\n\t# In case not shutdown cleanly.\n\t# shellcheck disable=SC2154\n\t[ -n \"$init\" ] && rm -f /etc/dfs/sharetab\n\n\t# Just simplify code later on.\n\tif [ -n \"$USE_DISK_BY_ID\" ] && [ \"$USE_DISK_BY_ID\" != 'yes' ]\n\tthen\n\t\t# It's something, but not 'yes' so it's no good to us.\n\t\tunset USE_DISK_BY_ID\n\tfi\n\n\t# Find list of already imported pools.\n\talready_imported=$(find_pools \"$ZPOOL\" list -H -oname)\n\tavailable_pools=$(find_pools \"$ZPOOL\" import)\n\n\t# Just in case - seen it happen (that a pool isn't visible/found\n\t# with a simple \"zpool import\" but only when using the \"-d\"\n\t# option or setting ZPOOL_IMPORT_PATH).\n\tif [ -d \"/dev/disk/by-id\" ]\n\tthen\n\t\tnpools=$(find_pools \"$ZPOOL\" import -d /dev/disk/by-id)\n\t\tif [ -n \"$npools\" ]\n\t\tthen\n\t\t\t# Because we have found extra pool(s) here, which wasn't\n\t\t\t# found 'normally', we need to force USE_DISK_BY_ID to\n\t\t\t# make sure we're able to actually import it/them later.\n\t\t\tUSE_DISK_BY_ID='yes'\n\n\t\t\tif [ -n \"$available_pools\" ]\n\t\t\tthen\n\t\t\t\t# Filter out duplicates (pools found with the simpl\n\t\t\t\t# \"zpool import\" but which is also found with the\n\t\t\t\t# \"zpool import -d ...\").\n\t\t\t\tnpools=$(echo \"$npools\" | sed \"s,$available_pools,,\")\n\n\t\t\t\t# Add the list to the existing list of\n\t\t\t\t# available pools\n\t\t\t\tavailable_pools=\"$available_pools;$npools\"\n\t\t\telse\n\t\t\t\tavailable_pools=\"$npools\"\n\t\t\tfi\n\t\tfi\n\tfi\n\n\t# Filter out any exceptions...\n\tif [ -n \"$ZFS_POOL_EXCEPTIONS\" ]\n\tthen\n\t\tlocal found=\"\"\n\t\tlocal apools=\"\"\n\t\tOLD_IFS=\"$IFS\" ; IFS=\";\"\n\n\t\tfor pool in $available_pools\n\t\tdo\n\t\t\tfor exception in $ZFS_POOL_EXCEPTIONS\n\t\t\tdo\n\t\t\t\t[ \"$pool\" = \"$exception\" ] && continue 2\n\t\t\t\tfound=\"$pool\"\n\t\t\tdone\n\n\t\t\tif [ -n \"$found\" ]\n\t\t\tthen\n\t\t\t\tif [ -n \"$apools\" ]\n\t\t\t\tthen\n\t\t\t\t\tapools=\"$apools;$pool\"\n\t\t\t\telse\n\t\t\t\t\tapools=\"$pool\"\n\t\t\t\tfi\n\t\t\tfi\n\t\tdone\n\n\t\tIFS=\"$OLD_IFS\"\n\t\tavailable_pools=\"$apools\"\n\tfi\n\n\t# For backwards compatibility, make sure that ZPOOL_IMPORT_PATH is set\n\t# to something we can use later with the real import(s). We want to\n\t# make sure we find all by* dirs, BUT by-vdev should be first (if it\n\t# exists).\n\tif [ -n \"$USE_DISK_BY_ID\" ] && [ -z \"$ZPOOL_IMPORT_PATH\" ]\n\tthen\n\t\tlocal dirs\n\t\tdirs=\"$(for dir in $(echo /dev/disk/by-*)\n\t\tdo\n\t\t\t# Ignore by-vdev here - we want it first!\n\t\t\techo \"$dir\" | grep -q /by-vdev && continue\n\t\t\t[ ! -d \"$dir\" ] && continue\n\n\t\t\tprintf \"%s\" \"$dir:\"\n\t\tdone | sed 's,:$,,g')\"\n\n\t\tif [ -d \"/dev/disk/by-vdev\" ]\n\t\tthen\n\t\t\t# Add by-vdev at the beginning.\n\t\t\tZPOOL_IMPORT_PATH=\"/dev/disk/by-vdev:\"\n\t\tfi\n\n\t\t# Help with getting LUKS partitions etc imported.\n\t\tif [ -d \"/dev/mapper\" ]; then\n\t\t\tif [ -n \"$ZPOOL_IMPORT_PATH\" ]; then\n\t\t\t\tZPOOL_IMPORT_PATH=\"$ZPOOL_IMPORT_PATH:/dev/mapper:\"\n\t\t\telse\n\t\t\t\tZPOOL_IMPORT_PATH=\"/dev/mapper:\"\n\t\t\tfi\n\t\tfi\n\n\t\t# ... and /dev at the very end, just for good measure.\n\t\tZPOOL_IMPORT_PATH=\"$ZPOOL_IMPORT_PATH$dirs:/dev\"\n\tfi\n\n\t# Needs to be exported for \"zpool\" to catch it.\n\t[ -n \"$ZPOOL_IMPORT_PATH\" ] && export ZPOOL_IMPORT_PATH\n\n\t# Mount all available pools (except those set in ZFS_POOL_EXCEPTIONS.\n\t#\n\t# If not interactive (run from init - variable init='/sbin/init')\n\t# we get ONE line for all pools being imported, with just a dot\n\t# as status for each pool.\n\t# Example: Importing ZFS pool(s)...                             [OK]\n\t#\n\t# If it IS interactive (started from the shell manually), then we\n\t# get one line per pool importing.\n\t# Example: Importing ZFS pool pool1                             [OK]\n\t#          Importing ZFS pool pool2                             [OK]\n\t#          [etc]\n\t[ -n \"$init\" ] && zfs_log_begin_msg \"Importing ZFS pool(s)\"\n\tOLD_IFS=\"$IFS\" ; IFS=\";\"\n\tfor pool in $available_pools\n\tdo\n\t\t[ -z \"$pool\" ] && continue\n\n\t\t# We have pools that haven't been imported - import them\n\t\tif [ -n \"$init\" ]\n\t\tthen\n\t\t\t# Not interactive - a dot for each pool.\n\t\t\t# Except on Gentoo where this doesn't work.\n\t\t\tzfs_log_progress_msg \".\"\n\t\telse\n\t\t\t# Interactive - one 'Importing ...' line per pool\n\t\t\tzfs_log_begin_msg \"Importing ZFS pool $pool\"\n\t\tfi\n\n\t\t# Import by using ZPOOL_IMPORT_PATH (either set above or in\n\t\t# the config file) _or_ with the 'built in' default search\n\t\t# paths. This is the preferred way.\n\t\t# shellcheck disable=SC2086\n\t\t\"$ZPOOL\" import -N ${ZPOOL_IMPORT_OPTS} \"$pool\" 2> /dev/null\n\t\tr=\"$?\" ; RET=$((RET + r))\n\t\tif [ \"$r\" -eq 0 ]\n\t\tthen\n\t\t\t# Output success and process the next pool\n\t\t\t[ -z \"$init\" ] && zfs_log_end_msg 0\n\t\t\tcontinue\n\t\tfi\n\t\t# We don't want a fail msg here, we're going to try import\n\t\t# using the cache file soon and that might succeed.\n\t\t[ ! -f \"$ZPOOL_CACHE\" ] && zfs_log_end_msg \"$RET\"\n\n\t\tif [ \"$r\" -gt 0 ] && [ -f \"$ZPOOL_CACHE\" ]\n\t\tthen\n\t\t\t# Failed to import without a cache file. Try WITH...\n\t\t\tif [ -z \"$init\" ] && check_boolean \"$VERBOSE_MOUNT\"\n\t\t\tthen\n\t\t\t\t# Interactive + Verbose = more information\n\t\t\t\tzfs_log_progress_msg \" using cache file\"\n\t\t\tfi\n\n\t\t\t# shellcheck disable=SC2086\n\t\t\t\"$ZPOOL\" import -c \"$ZPOOL_CACHE\" -N ${ZPOOL_IMPORT_OPTS} \\\n\t\t\t\t\"$pool\" 2> /dev/null\n\t\t\tr=\"$?\" ; RET=$((RET + r))\n\t\t\tif [ \"$r\" -eq 0 ]\n\t\t\tthen\n\t\t\t\t[ -z \"$init\" ] && zfs_log_end_msg 0\n\t\t\t\tcontinue 3 # Next pool\n\t\t\tfi\n\t\t\tzfs_log_end_msg \"$RET\"\n\t\tfi\n\tdone\n\t[ -n \"$init\" ] && zfs_log_end_msg \"$RET\"\n\n\tIFS=\"$OLD_IFS\"\n\t[ -n \"$already_imported\" ] && [ -z \"$available_pools\" ] && return 0\n\n\treturn \"$RET\"\n}\n\ndo_import()\n{\n\tif check_boolean \"$ZPOOL_IMPORT_ALL_VISIBLE\"\n\tthen\n\t\tdo_import_all_visible\n\telse\n\t\t# This is the default option\n\t\tdo_verbatim_import\n\tfi\n}\n\n# Output the status and list of pools\ndo_status()\n{\n\tcheck_module_loaded \"zfs\" || exit 0\n\n\t\"$ZPOOL\" status && echo \"\" && \"$ZPOOL\" list\n}\n\ndo_start()\n{\n\tif check_boolean \"$VERBOSE_MOUNT\"\n\tthen\n\t    zfs_log_begin_msg \"Checking if ZFS userspace tools present\"\n\tfi\n\n\tif checksystem\n\tthen\n\t\tcheck_boolean \"$VERBOSE_MOUNT\" && zfs_log_end_msg 0\n\n\t\tcheck_boolean \"$VERBOSE_MOUNT\" && \\\n\t\t\tzfs_log_begin_msg \"Loading kernel ZFS infrastructure\"\n\n\t\tif ! load_module \"zfs\"\n\t\tthen\n\t\t\tcheck_boolean \"$VERBOSE_MOUNT\" && zfs_log_end_msg 1\n\t\t\treturn 5\n\t\tfi\n\t\tcheck_boolean \"$VERBOSE_MOUNT\" && zfs_log_end_msg 0\n\n\t\tdo_import && udev_trigger # just to make sure we get zvols.\n\n\t\treturn 0\n\telse\n\t\treturn 1\n\tfi\n}\n\n# ----------------------------------------------------\n\nif [ ! -e /sbin/openrc-run ]\nthen\n\tcase \"$1\" in\n\t\tstart)\n\t\t\tdo_start\n\t\t\t;;\n\t\tstop)\n\t\t\t# no-op\n\t\t\t;;\n\t\tstatus)\n\t\t\tdo_status\n\t\t\t;;\n\t\tforce-reload|condrestart|reload|restart)\n\t\t\t# no-op\n\t\t\t;;\n\t\t*)\n\t\t\t[ -n \"$1\" ] && echo \"Error: Unknown command $1.\"\n\t\t\techo \"Usage: $0 {start|status}\"\n\t\t\texit 3\n\t\t\t;;\n\tesac\n\n\texit $?\nelse\n\t# Create wrapper functions since Gentoo don't use the case part.\n\tdepend() { do_depend; }\n\tstart() { do_start; }\n\tstatus() { do_status; }\nfi\n",
  "logic_map": {},
  "failure_modes": [],
  "crash_correlation_map": {}
}