# oss_sovereignty/zfs-2.2.2/etc/init.d/zfs-load-key.in.mpy
{
  "module_name": "zfs-load-key.in",
  "hash_id": "f6706787d42db271544443bb51354bdf6428c5637ceb32db52b75c193d1bf225",
  "original_prompt": "Ingested from zfs-2.2.2/etc/init.d/zfs-load-key.in",
  "human_readable_source": "#!@DEFAULT_INIT_SHELL@\n# shellcheck disable=SC2154\n#\n# zfs-load-key  This script will load/unload the zfs filesystems keys.\n#\n# chkconfig:    2345 06 99\n# description:  This script will load or unload the zfs filesystems keys during\n#               system boot/shutdown. Only filesystems with key path set\n#               in keylocation property. See the zfs(8) man page for details.\n# probe: true\n#\n### BEGIN INIT INFO\n# Provides:          zfs-load-key\n# Required-Start:    $local_fs zfs-import\n# Required-Stop:     $local_fs zfs-import\n# Default-Start:     2 3 4 5\n# Default-Stop:      0 1 6\n# X-Start-Before:    zfs-mount\n# X-Stop-After:      zfs-zed\n# Short-Description: Load ZFS keys for filesystems and volumes\n# Description: Run the `zfs load-key` or `zfs unload-key` commands.\n### END INIT INFO\n#\n# Released under the 2-clause BSD license.\n#\n# This script is based on debian/zfsutils.zfs.init from the\n# Debian GNU/kFreeBSD zfsutils 8.1-3 package, written by Aurelien Jarno.\n\n# Source the common init script\n. @sysconfdir@/zfs/zfs-functions\n\n# ----------------------------------------------------\n\ndo_depend()\n{\n\t# bootmisc will log to /var which may be a different zfs than root.\n\tbefore bootmisc logger zfs-mount\n\n\tafter zfs-import sysfs\n\tkeyword -lxc -openvz -prefix -vserver\n}\n\n# Load keys for all datasets/filesystems\ndo_load_keys()\n{\n\tzfs_log_begin_msg \"Load ZFS filesystem(s) keys\"\n\n\t\"$ZFS\" list -Ho name,encryptionroot,keystatus,keylocation |\n\t    while IFS=\"\t\" read -r name encryptionroot keystatus keylocation; do\n\t\tif [ \"$encryptionroot\" != \"-\" ] &&\n\t\t\t[ \"$name\" = \"$encryptionroot\" ] &&\n\t\t\t[ \"$keystatus\" = \"unavailable\" ] &&\n\t\t\t[ \"$keylocation\" != \"prompt\" ] &&\n\t\t\t[ \"$keylocation\" != \"none\" ]\n\t\tthen\n\t\t\tzfs_action \"Load key for $encryptionroot\" \\\n\t\t\t    \"$ZFS\" load-key \"$encryptionroot\"\n\t\tfi\n\tdone\n\n\tzfs_log_end_msg 0\n\n\treturn 0\n}\n\n# Unload keys for all datasets/filesystems\ndo_unload_keys()\n{\n\tzfs_log_begin_msg \"Unload ZFS filesystem(s) key\"\n\n\t\"$ZFS\" list -Ho name,encryptionroot,keystatus | sed '1!G;h;$!d' |\n\t    while IFS=\"\t\" read -r name encryptionroot keystatus; do\n\t\tif [ \"$encryptionroot\" != \"-\" ] &&\n\t\t\t[ \"$name\" = \"$encryptionroot\" ] &&\n\t\t\t[ \"$keystatus\" = \"available\" ]\n\t\tthen\n\t\t\tzfs_action \"Unload key for $encryptionroot\" \\\n\t\t\t    \"$ZFS\" unload-key \"$encryptionroot\"\n\t\tfi\n\tdone\n\n\tzfs_log_end_msg 0\n\n\treturn 0\n}\n\ndo_start()\n{\n\tcheck_boolean \"$ZFS_LOAD_KEY\" || exit 0\n\n\tcheck_module_loaded \"zfs\" || exit 0\n\n\tdo_load_keys\n}\n\ndo_stop()\n{\n\tcheck_boolean \"$ZFS_UNLOAD_KEY\" || exit 0\n\n\tcheck_module_loaded \"zfs\" || exit 0\n\n\tdo_unload_keys\n}\n\n# ----------------------------------------------------\n\nif [ ! -e /sbin/openrc-run ]\nthen\n\tcase \"$1\" in\n\t\tstart)\n\t\t\tdo_start\n\t\t\t;;\n\t\tstop)\n\t\t\tdo_stop\n\t\t\t;;\n\t\tforce-reload|condrestart|reload|restart|status)\n\t\t\t# no-op\n\t\t\t;;\n\t\t*)\n\t\t\t[ -n \"$1\" ] && echo \"Error: Unknown command $1.\"\n\t\t\techo \"Usage: $0 {start|stop}\"\n\t\t\texit 3\n\t\t\t;;\n\tesac\n\n\texit $?\nelse\n\t# Create wrapper functions since Gentoo don't use the case part.\n\tdepend() { do_depend; }\n\tstart() { do_start; }\n\tstop() { do_stop; }\nfi\n",
  "logic_map": {},
  "failure_modes": [],
  "crash_correlation_map": {}
}