# oss_sovereignty/zfs-2.2.2/etc/init.d/zfs-mount.in.mpy
{
  "module_name": "zfs-mount.in",
  "hash_id": "9257a3162db837a95e6fca4bbe52621df448b9d34e2cc8254a47a0a6910f9bbc",
  "original_prompt": "Ingested from zfs-2.2.2/etc/init.d/zfs-mount.in",
  "human_readable_source": "#!@DEFAULT_INIT_SHELL@\n# shellcheck disable=SC2154\n#\n# zfs-mount     This script will mount/umount the zfs filesystems.\n#\n# chkconfig:    2345 06 99\n# description:  This script will mount/umount the zfs filesystems during\n#               system boot/shutdown. Configuration of which filesystems\n#               should be mounted is handled by the zfs 'mountpoint' and\n#               'canmount' properties. See the zfs(8) man page for details.\n#               It is also responsible for all userspace zfs services.\n# probe: true\n#\n### BEGIN INIT INFO\n# Provides:          zfs-mount\n# Required-Start:    zfs-import\n# Required-Stop:     $local_fs zfs-import\n# Default-Start:     S\n# Default-Stop:      0 1 6\n# X-Start-Before:    mountall\n# X-Stop-After:      zfs-zed\n# Short-Description: Mount ZFS filesystems and volumes\n# Description: Run the `zfs mount -a` or `zfs umount -a` commands.\n### END INIT INFO\n#\n# Released under the 2-clause BSD license.\n#\n# This script is based on debian/zfsutils.zfs.init from the\n# Debian GNU/kFreeBSD zfsutils 8.1-3 package, written by Aurelien Jarno.\n\n# Source the common init script\n. @sysconfdir@/zfs/zfs-functions\n\n# ----------------------------------------------------\n\nchkroot() {\n\twhile read -r _ mp _; do\n\t\tif [ \"$mp\" = \"/\" ]; then\n\t\t\treturn 0\n\t\tfi\n\tdone < /proc/self/mounts\n\n\treturn 1\n}\n\ndo_depend()\n{\n\t# Try to allow people to mix and match fstab with ZFS in a way that makes sense.\n\tif [ \"$(mountinfo -s /)\" = 'zfs' ]\n\tthen\n\t\tbefore localmount\n\telse\n\t\tafter localmount\n\tfi\n\n\t# bootmisc will log to /var which may be a different zfs than root.\n\tbefore bootmisc logger\n\n\tafter zfs-import sysfs\n\tuse mtab\n\tkeyword -lxc -openvz -prefix -vserver\n}\n\n# Mount all datasets/filesystems\ndo_mount()\n{\n\tlocal verbose overlay\n\n\tcheck_boolean \"$VERBOSE_MOUNT\" && verbose=v\n\tcheck_boolean \"$DO_OVERLAY_MOUNTS\" && overlay=O\n\n\tzfs_action \"Mounting ZFS filesystem(s)\" \\\n\t    \"$ZFS\" mount \"-a$verbose$overlay\" \"$MOUNT_EXTRA_OPTIONS\"\n\n\treturn 0\n}\n\n# Unmount all filesystems\ndo_unmount()\n{\n\t# This shouldn't really be necessary, as long as one can get\n\t# zfs-import to run sufficiently late in the shutdown/reboot process\n\t# - after unmounting local filesystems. This is just here in case/if\n\t# this isn't possible.\n\tzfs_action \"Unmounting ZFS filesystems\" \"$ZFS\" unmount -a\n\n\treturn 0\n}\n\ndo_start()\n{\n\tcheck_boolean \"$ZFS_MOUNT\" || exit 0\n\n\tcheck_module_loaded \"zfs\" || exit 0\n\n\t# Ensure / exists in /proc/self/mounts.\n\t# This should be handled by rc.sysinit but lets be paranoid.\n\tif ! chkroot\n\tthen\n\t\tmount -f /\n\tfi\n\n\tdo_mount\n}\n\ndo_stop()\n{\n\tcheck_boolean \"$ZFS_UNMOUNT\" || exit 0\n\n\tcheck_module_loaded \"zfs\" || exit 0\n\n\tdo_unmount\n}\n\n# ----------------------------------------------------\n\nif [ ! -e /sbin/openrc-run ]\nthen\n\tcase \"$1\" in\n\t\tstart)\n\t\t\tdo_start\n\t\t\t;;\n\t\tstop)\n\t\t\tdo_stop\n\t\t\t;;\n\t\tforce-reload|condrestart|reload|restart|status)\n\t\t\t# no-op\n\t\t\t;;\n\t\t*)\n\t\t\t[ -n \"$1\" ] && echo \"Error: Unknown command $1.\"\n\t\t\techo \"Usage: $0 {start|stop}\"\n\t\t\texit 3\n\t\t\t;;\n\tesac\n\n\texit $?\nelse\n\t# Create wrapper functions since Gentoo don't use the case part.\n\tdepend() { do_depend; }\n\tstart() { do_start; }\n\tstop() { do_stop; }\nfi\n",
  "logic_map": {},
  "failure_modes": [],
  "crash_correlation_map": {}
}